# upstart.conf
#!upstart
# Fichier de configuration upstart
description "moi-candidat"

start on started mountall
stop on shutdown

# Automatically Respawn:
respawn
respawn limit 99 5

env NODE_ENV=production

exec "/root/projects/moi-candidat/run.sh"

